# rootfs/etc/supervisor.d/supervisord.ini
[supervisord]
nodaemon = true

[program:squid]
command=/usr/sbin/squid -N -f /etc/squid/squid.conf
redirect_stderr=true
stdout_logfile=/var/log/supervisor/squid.info.log
stderr_logfile=/var/log/supervisor/squid.error.log
stdout_logfile_maxbytes=16MB
stderr_logfile_maxbytes=16MB

[program:nghttpx]
command=nghttpx
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nghttpx.info.log
stderr_logfile=/var/log/supervisor/nghttpx.error.log
stdout_logfile_maxbytes=16MB
stderr_logfile_maxbytes=16MB

